# deploy/systemd/binance-sync-docker.timer
[Unit]
Description=Daily Binance Futures Data Sync (Docker)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
